# init.d_scripts/mcp-crawl4ai-rag.service
[Unit]
Description=MCP Crawl4AI RAG Server
Documentation=https://github.com/your-repo/mcp-crawl4ai-rag
After=network.target docker.service
Requires=docker.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=aitoolset
Group=aitoolset
WorkingDirectory=/home/aitoolset/mcp-crawl4ai-rag
Environment=PATH=/home/aitoolset/mcp-crawl4ai-rag/.venv/bin:/home/aitoolset/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/aitoolset/mcp-crawl4ai-rag
EnvironmentFile=/home/aitoolset/mcp-crawl4ai-rag/.env
ExecStartPre=/home/aitoolset/mcp-crawl4ai-rag/.venv/bin/python /home/aitoolset/mcp-crawl4ai-rag/setup_local_supabase.py
ExecStart=/home/aitoolset/mcp-crawl4ai-rag/.venv/bin/python /home/aitoolset/mcp-crawl4ai-rag/src/crawl4ai_mcp.py
ExecStop=/usr/bin/docker compose -f /home/aitoolset/mcp-crawl4ai-rag/docker-compose.yml down
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-crawl4ai-rag

# Reduce security restrictions to allow Docker access
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# Resource limits
LimitNOFILE=65536
MemoryMax=8G

[Install]
WantedBy=multi-user.target
